# pkg/console/debootstrap/conf/lxde_a/pattern/rc.local
#!/bin/sh -e

# VladVons
# force DNS
echo "nameserver 8.8.8.8" > /etc/resolv.conf

exit 0
